# launchd/com.opencode.parakeet-stt.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Template: setup.sh replaces HOME_PLACEHOLDER and writes the selected accelerator. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.opencode.parakeet-stt</string>
  <key>ProgramArguments</key>
  <array>
    <string>HOME_PLACEHOLDER/.config/opencode/parakeet-stt/.venv/bin/python</string>
    <string>HOME_PLACEHOLDER/.config/opencode/parakeet-stt/server.py</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key><string>HOME_PLACEHOLDER</string>
    <key>PATH</key><string>HOME_PLACEHOLDER/.config/opencode/parakeet-stt/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>PARAKEET_PORT</key><string>5093</string>
    <key>PARAKEET_USE_GPU</key><string>false</string>
    <key>PARAKEET_DEFAULT_MODEL</key><string>parakeet-tdt-0.6b-v3</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>WorkingDirectory</key><string>HOME_PLACEHOLDER/.config/opencode/parakeet-stt</string>
  <key>StandardOutPath</key><string>HOME_PLACEHOLDER/.config/opencode/parakeet-stt.log</string>
  <key>StandardErrorPath</key><string>HOME_PLACEHOLDER/.config/opencode/parakeet-stt.log</string>
</dict>
</plist>
