# launchd/com.opencode.supertonic.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Template: setup.sh replaces HOME_PLACEHOLDER and writes the selected port/backend. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.opencode.supertonic</string>
  <key>ProgramArguments</key>
  <array>
    <string>HOME_PLACEHOLDER/.config/opencode/supertonic-tts/.venv/bin/python</string>
    <string>-m</string>
    <string>uvicorn</string>
    <string>api.src.main:app</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8766</string>
    <string>--app-dir</string>
    <string>HOME_PLACEHOLDER/.config/opencode/supertonic-tts/py</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>HOME_PLACEHOLDER</string>
    <key>PATH</key>
    <string>HOME_PLACEHOLDER/.config/opencode/supertonic-tts/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>SUPERTONIC_MODEL_DIR</key>
    <string>HOME_PLACEHOLDER/.config/opencode/supertonic-tts/assets/supertonic-3</string>
    <key>ONNX_DIR</key>
    <string>HOME_PLACEHOLDER/.config/opencode/supertonic-tts/assets/supertonic-3/onnx</string>
    <key>VOICE_STYLES_DIR</key>
    <string>HOME_PLACEHOLDER/.config/opencode/supertonic-tts/assets/supertonic-3/voice_styles</string>
    <key>SUPERTONIC_MLX_MODEL_DIR</key>
    <string>HOME_PLACEHOLDER/.config/opencode/supertonic-tts/assets/supertonic-3-mlx</string>
    <key>SUPERTONIC_MLX_AUTO_DOWNLOAD</key>
    <string>false</string>
    <key>SUPERTONIC_MLX_FALLBACK_TO_ONNX</key>
    <string>true</string>
    <key>USE_GPU</key>
    <string>false</string>
    <key>SUPERTONIC_ORT_BACKEND</key>
    <string>auto</string>
    <key>LOG_LEVEL</key>
    <string>INFO</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>WorkingDirectory</key>
  <string>HOME_PLACEHOLDER/.config/opencode/supertonic-tts/py</string>
  <key>StandardOutPath</key>
  <string>HOME_PLACEHOLDER/.config/opencode/supertonic.log</string>
  <key>StandardErrorPath</key>
  <string>HOME_PLACEHOLDER/.config/opencode/supertonic.log</string>
</dict>
</plist>
